# gitblog/nginx.ini
[program:nginx]
command=/usr/local/nginx/sbin/nginx -g "daemon off;"
directory=/usr/local/nginx/
autostart=true
autorestart=unexpected
startsecs=1
user=root