# gitblog/php-fpm.ini
[program:php-fpm]
command=/usr/local/php/sbin/php-fpm -F
directory=/usr/local/php/
autostart=true
autorestart=unexpected
startsecs=1
user=root